# setup/pi-config/tourguide-ws.service
[Unit]
Description=Tour Guide WebSocket Audio Server
After=network.target

[Service]
Type=simple
User=pi
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /usr/local/bin/tourguide-ws-server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
